# presentation/srv/chat_upstart.conf
#!upstart
description "Tornado chat app"

start on startup
stop on shutdown
respawn

setgid dev
setuid dev

chdir /home/dev/tornado/demos/chat
exec /home/dev/venv/bin/python chatdemo.py
